# redisio/templates/default/redis@.service.erb
[Unit]
Description=Redis (<%= @server_name %>) persistent key-value database
After=network.target

[Service]
ExecStart=<%= @bin_path %>/redis-server /etc/redis/<%= @server_name %>.conf --daemonize no
User=<%= @user %>
Group=<%= @group %>
LimitNOFILE=<%= @limit_nofile %>

[Install]
WantedBy=multi-user.target
