# base/templates/default/dnsmasq.conf.erb
[program:dnsmasq]
command=/usr/sbin/dnsmasq -d -k -u dnsmasq -7 /etc/dnsmasq.d --local-service
autostart=true
autorestart=true
stopasgroup=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_err.log