# base/templates/default/php-fpm.conf.erb
[program:php<%= node["php"]["version"] %>-fpm]
command=/bin/bash -c 'source /etc/environment && /usr/sbin/php-fpm<%= node["php"]["version"] %> --nodaemonize --fpm-config /etc/php/<%= node["php"]["version"] %>/fpm/php-fpm.conf'
autostart=false
autorestart=unexpected
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_err.log
exitcodes=0